# .config/autostart/open-journal.desktop
[Desktop Entry]
Type=Application
Exec=fish -c 'subl "$JOURNAL_PATH"'
X-GNOME-Autostart-enabled=true
Name=Open Journal
