# .config/autostart/open-journal.desktop
[Desktop Entry]
Type=Application
Exec=sh -c "subl \"\$HOME/Work/2018-CeleraOne/journal.md\""
X-GNOME-Autostart-enabled=true
Name=Open Journal
